# tunetailorbot.service
[Unit]
Description=TuneTailor Discord Bot
After=network.target redis6.service
Wants=redis6.service

[Service]
Type=simple
User=ec2-user
WorkingDirectory=/home/ec2-user/TuneTailor-Deployment-220
ExecStart=/usr/bin/java -jar /home/ec2-user/TuneTailor-Deployment-220/target/*jar-with-dependencies.jar
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target